# installation/minecraft.service
[Unit]
Description=Minecraft server
After=network.target

[Service]
Type=forking
ExecStart=/bin/sh /usr/local/bin/minecraft/start.sh
#ExecStop=/bin/sh /usr/local/bin/minecraft/stop.sh
TimeoutSec=600

[Install]
WantedBy=multi-user.target
